# deploy/hetzner/cron/qdrant-backup.cron
# Weekly Qdrant snapshot backup (every Sunday at 3 AM)
# Install with: sudo cp deploy/hetzner/cron/qdrant-backup.cron /etc/cron.d/qdrant-backup
SHELL=/bin/bash
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
0 3 * * 0 root cd /opt/RAG-pipeline && /bin/bash -lc 'set -a && source .env && set +a && ./deploy/hetzner/backup.sh' >> /var/log/qdrant-backup.log 2>&1
